# files/lldp2facts.timer
[Unit]
Description=run lldp renew 15min

[Timer]
OnUnitActiveSec=15min
OnBootSec=10min
Persistent=true

[Install]
WantedBy=multi-user.target
