# files/lldp2facts.service
[Unit]
Description=Run lldpctl and update facts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lldp2facts
