# files/lldp2facts.cron
*/15 * * * * root /usr/local/bin/lldp2facts
